# profile.d/40-environment
#!/bin/bash

export GEM_HOME="$HOME/.gems"
PATH="$HOME/.gems/bin:$PATH"
